# sys_monitor_http/installation_and_config/systemd_cfg/sysmon.socket
[Unit]
Description=Sys Monitor HTTP - Socket

[Socket]
ListenStream=16801

[Install]
WantedBy=sockets.target
